# spectrum/image/rootfs/usr/lib/systemd/user/org.gnome.Shell@.service.d/70-spectrum-software-cursor.conf
[Service]
# NVIDIA's hardware cursor plane can disappear when a Wayland surface enters
# fullscreen even though pointer input continues to work. Composite the cursor
# with the rest of the GNOME Shell scene until the driver bug is fixed.
Environment=MUTTER_DEBUG_DISABLE_HW_CURSORS=1
