# ansible/files/scripts/host/keyboard/toshy-kanata-device.service
[Unit]
Description=Refresh Toshy after Kanata virtual keyboard changes
Documentation=https://github.com/RedBearAK/Toshy https://github.com/jtroo/kanata

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'systemctl --user reset-failed toshy-config.service || true; [ -e /run/kanata-main/main ] || exit 0; if systemctl --user is-active --quiet toshy-config.service || systemctl --user is-enabled --quiet toshy-config.service; then systemctl --user restart toshy-config.service; fi'
