# spectrum/image/rootfs/usr/lib/systemd/user/uresourced.service.d/99-spectrum-disable-user-daemon.conf
[Unit]
# Keep the user-session uresourced daemon opt-in on Spectrum.
#
# Why:
# - On Bluefin/Fedora 44 with uresourced-0.5.4-5.fc44 and systemd-259.7,
#   /usr/libexec/uresourced --user was observed looping on user-manager
#   SetUnitProperties calls for GNOME/session units. In that state, new user
#   systemd transient jobs remained forever in "start waiting".
# - A trivial `systemd-run --user --scope /usr/bin/true` hung, and Flatpak
#   sandbox commands never started. Killing/masking only the user uresourced
#   daemon made both commands return immediately and Flatpak apps launched
#   again.
# - Fedora documents uresourced as a cgroup/resource-protection layer for the
#   active graphical session. It also notes that the daemon was designed to be
#   obsoleted once the behavior moves into other upstream components. Spectrum
#   already sets explicit resource controls for important desktop units, so the
#   risk of losing this dynamic user daemon is lower than the risk of wedging
#   all transient app launches.
#
# This does not remove the uresourced package or the system daemon. It only
# prevents the GNOME/user-session daemon from starting by default. To test it
# again in a future Fedora/systemd build, create this runtime marker before the
# graphical session starts:
#
#   sudo mkdir -p /run/spectrum
#   sudo touch /run/spectrum/enable-user-uresourced
#
# Useful regression checks:
#
#   timeout 5s systemd-run --user --scope /usr/bin/true
#   timeout 5s flatpak run --command=sh org.flatpak.Builder -lc 'echo ok'
#
# Both should complete without leaving jobs in `systemctl --user list-jobs`.
#
# References:
# - Fedora package summary:
#   https://packages.fedoraproject.org/pkgs/uresourced/uresourced/
# - Fedora Workstation change that introduced uresourced:
#   https://fedoraproject.org/wiki/Changes/Reserve_resources_for_active_user_WS
# - systemd-run transient unit behavior:
#   https://www.freedesktop.org/software/systemd/man/latest/systemd-run.html
# - Flatpak command reference:
#   https://docs.flatpak.org/en/latest/flatpak-command-reference.html
ConditionPathExists=/run/spectrum/enable-user-uresourced
