# dotfiles/dot_config/systemd/user/carmenta-prewarm.service
[Unit]
Description=Prewarm Carmenta emoji picker
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/carmenta-activate --prewarm
CapabilityBoundingSet=
LockPersonality=true
NoNewPrivileges=true
PrivateDevices=true
PrivateNetwork=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=default.target
