# dotfiles/dot_config/systemd/user/solaar.service
[Unit]
Description=Solaar Logitech device manager
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/run/current-system/sw/bin
ExecStart=%h/.local/bin/solaar-no-tray
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
